# rootfs/etc/s6-overlay/s6-rc.d/retropanel/run
#!/usr/bin/with-contenv bashio
cd /usr/lib/retropanel
exec python3 server.py
