# rootfs/etc/s6-overlay/s6-rc.d/init-retropanel/up
#!/usr/bin/with-contenv bashio
bashio::log.info "Retro Panel starting..."
bashio::log.info "HA URL: $(bashio::config 'ha_url')"
bashio::log.info "Panel title: $(bashio::config 'panel_title')"
